# scripts/systemd/cidrella.service
[Unit]
Description=CIDRella IPAM Server
After=network.target cidrella-dnsmasq.service
Wants=cidrella-dnsmasq.service

[Service]
Type=simple
User=cidrella
Group=cidrella
WorkingDirectory=/opt/cidrella
ExecStart=/usr/bin/node server/src/index.js
Environment=DATA_DIR=/var/lib/cidrella
Environment=NODE_ENV=production
Environment=HTTPS_PORT=8443
Environment=HTTP_PORT=8080
Restart=on-failure
RestartSec=5

# Security hardening
ProtectSystem=strict
ReadWritePaths=/var/lib/cidrella
ProtectHome=true
PrivateTmp=true

# Network scanning requires sudo for arping/nmap (NoNewPrivileges must be off)
AmbientCapabilities=CAP_NET_RAW CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
