# rootfs/etc/s6-overlay/s6-rc.d/anomaly-detector/run
#!/command/execlineb -P
s6-setuidgid cidrella
python3 /app/server/anomaly/anomaly_detector.py
