# rootfs/etc/s6-overlay/s6-rc.d/node-app/run
#!/command/execlineb -P
exec node /app/server/src/index.js
