# rootfs/etc/s6-overlay/s6-rc.d/node-app/run
#!/command/execlineb -P
s6-setuidgid cidrella
node /app/server/src/index.js
