# rootfs/etc/s6-overlay/s6-rc.d/node-app/run
#!/command/execlineb -P
importas -D "0" PUID PUID
ifelse { test ${PUID} -ne 0 }
{ s6-setuidgid ipam node /app/server/src/index.js }
exec node /app/server/src/index.js
