# rootfs/etc/s6-overlay/s6-rc.d/dnsmasq/run
#!/command/execlineb -P
exec dnsmasq --keep-in-foreground --conf-file=/data/dnsmasq/dnsmasq.conf
